# ansible/playbooks/templates/infrahub-agent.service.j2
[Unit]
Description=InfraHub Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/infrahub-agent
# URL master только из API (extra-var infrahub_api_base); argv переопределяет env и дефолт clap.
# Имя агента — короткое имя хоста с целевой ноды (gather_facts → ansible_hostname).
Environment="INFRAHUB_AGENT_NAME={{ ansible_hostname }}"
WorkingDirectory={{ infrahub_agent_state_dir }}
ExecStart={{ infrahub_agent_install_path }} --server={{ infrahub_api_base | trim }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
